# mcp-servers/cqp-search/cqp-search.service
[Unit]
Description=CQP Corpus Query Engine
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cqp-search
ExecStart=/usr/bin/python3 -m gunicorn -w 2 -b 0.0.0.0:5001 --timeout 60 server:app
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
